# com.wavetree.news-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
DEPRECATED (2026-02-21)
- LaunchAgent 기반 운영 중단: 외장 볼륨 경로 접근 시 Operation not permitted 발생 가능
- 현재 운영 기준: wave-tree-news-hub cron 파이프라인(06:50/07:00)
- 참고 문서: /Volumes/AI_DATA_CENTRE/AI_WORKSPACE/wave-tree-news-hub/AUTOMATION_STATUS.md
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.wavetree.news-sync</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>/Users/seunghoonoh/.gemini/antigravity/scratch/node-v20/bin/node /Volumes/AI_DATA_CENTRE/AI_WORKSPACE/wave-tree-news-hub/scripts/normalize.js --in /Volumes/AI_DATA_CENTRE/AI_WORKSPACE/wave-tree-news-hub/data/raw/perplexity.txt --out /Volumes/AI_DATA_CENTRE/AI_WORKSPACE/wave-tree-news-hub/data/normalized/news.json &amp;&amp; cp /Volumes/AI_DATA_CENTRE/AI_WORKSPACE/wave-tree-news-hub/data/normalized/news.json /Volumes/AI_DATA_CENTRE/AI_WORKSPACE/woonmok.github.io/news.json &amp;&amp; /Library/Frameworks/Python.framework/Versions/3.14/bin/python3 /Volumes/AI_DATA_CENTRE/AI_WORKSPACE/wave-tree-news-hub/sync_top_news.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Volumes/AI_DATA_CENTRE/AI_WORKSPACE/wave-tree-news-hub</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>StandardOutPath</key>
	<string>/tmp/com.wavetree.news-sync.out.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/com.wavetree.news-sync.err.log</string>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>300</integer>
</dict>
</plist>
