# com.wavetree.news-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.wavetree.news-sync</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Volumes/AI_DATA_CENTRE/AI_WORKSPACE/wave-tree-news-hub/.venv312/bin/python</string>
		<string>/Volumes/AI_DATA_CENTRE/AI_WORKSPACE/wave-tree-news-hub/sync_top_news.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Volumes/AI_DATA_CENTRE/AI_WORKSPACE/wave-tree-news-hub</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>StandardOutPath</key>
	<string>/Volumes/AI_DATA_CENTRE/AI_WORKSPACE/woonmok.github.io/logs/news_sync.log</string>
	<key>StandardErrorPath</key>
	<string>/Volumes/AI_DATA_CENTRE/AI_WORKSPACE/woonmok.github.io/logs/news_sync_error.log</string>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>300</integer>
</dict>
</plist>
