# com.wavetree.antigravity.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.wavetree.antigravity</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>cd /Volumes/AI_DATA_CENTRE/AI_WORKSPACE/woonmok.github.io &amp;&amp; exec /Library/Frameworks/Python.framework/Versions/3.14/bin/python3 antigravity.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/seunghoonoh</string>
	<key>StandardOutPath</key>
	<string>/tmp/com.wavetree.antigravity.out.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/com.wavetree.antigravity.err.log</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>RestartDelay</key>
	<integer>10</integer>
</dict>
</plist>
