# com.wavetree.news-sync-loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
DEPRECATED (2026-02-21)
- LaunchAgent 기반 운영 중단: 외장 볼륨 경로 접근 시 Operation not permitted 발생 가능
- 현재 운영 기준: wave-tree-news-hub cron 파이프라인(06:50/07:00)
- 참고 문서: ../wave-tree-news-hub/AUTOMATION_STATUS.md
-->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.wavetree.news-sync-loop</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>NEWS_HUB_DIR="${WAVETREE_NEWS_HUB_DIR:-${WAVETREE_WORKSPACE_ROOT:-$HOME/AI_WORKSPACE}/wave-tree-news-hub}"; exec /bin/bash "$NEWS_HUB_DIR/scripts/auto_sync_news_fswatch.sh"</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>StandardOutPath</key>
	<string>/tmp/com.wavetree.news-sync-loop.out.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/com.wavetree.news-sync-loop.err.log</string>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>300</integer>
</dict>
</plist>
