# com.wavetree.antigravity.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.wavetree.antigravity</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/seunghoonoh/woonmok.github.io/.venv/bin/python3</string>
		<string>/Users/seunghoonoh/woonmok.github.io/antigravity.py</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/seunghoonoh/woonmok.github.io</string>
	<key>StandardOutPath</key>
	<string>/Users/seunghoonoh/woonmok.github.io/logs/antigravity.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/seunghoonoh/woonmok.github.io/logs/antigravity_error.log</string>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>RestartDelay</key>
	<integer>10</integer>
</dict>
</plist>
